# emqtt.ini
[supervisord]
nodaemon=true

[program:emqtt]
command=ash -c ". /myemqenv; exec /opt/emqttd/start.sh"
priority=5
numprocs=1
autostart=true
autorestart=true
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:clusterize]
command=ash -c ". /myemqenv; exec /usr/local/bin/clusterize"
priority=10
numprocs=1
autostart=true
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
